# packages/stack/src/getbased_agent_stack/systemd/getbased-rag.service
[Unit]
Description=getbased-rag (local RAG knowledge server)
Documentation=https://github.com/elkimek/getbased-agents/tree/main/packages/rag
After=network.target

[Service]
Type=exec
# Shared env file written by `getbased-stack init` — carries GETBASED_TOKEN,
# LENS_API_KEY_FILE, LENS_URL (and the GETBASED_STACK_MANAGED=1 opt-in
# flag, which must be set for the Python env loader to pick up this file).
# Optional: /etc/default/getbased-rag for system-wide overrides.
EnvironmentFile=-%h/.config/getbased/env
EnvironmentFile=-/etc/default/getbased-rag
Environment=GETBASED_STACK_MANAGED=1
Environment=LENS_HOST=127.0.0.1
Environment=LENS_PORT=8322
ExecStart=%h/.local/bin/lens serve
# Restart=always (not on-failure) so a clean SIGTERM — e.g. operator stop
# during troubleshooting — still brings the service back at the next boot.
# Matches the failure mode that kept lens-rag.service dead for hours on
# the Hermes VM until someone noticed.
Restart=always
RestartSec=5s

# User-mode-compatible hardening only. Directives that require
# CAP_SYS_ADMIN (ProtectKernelTunables, ProtectKernelModules,
# ProtectControlGroups, RestrictNamespaces, RestrictAddressFamilies,
# RestrictRealtime, MemoryDenyWriteExecute) fail with 218/CAPABILITIES
# under `systemctl --user` on a non-privileged bus and are omitted.
# For a hardened system-service deployment, promote this unit to
# /etc/systemd/system/ and re-add those directives.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.local/share/getbased
PrivateTmp=true
LockPersonality=true

[Install]
WantedBy=default.target
